# platforms/debian/meshchatter.service
[Unit]
Description=Meshchatter
Wants=network-online.target
After=network-online.targe

[Service]
Type=simple
StandardOutput=journal
ExecStart=/root/meshchatter/main.uc

[Install]
WantedBy=default.target
