# platforms/debian/raven.service
[Unit]
Description=raven
Wants=network-online.target
After=network-online.targe

[Service]
Type=simple
StandardOutput=journal
ExecStart=/root/raven/main.uc

[Install]
WantedBy=default.target
